# deploy/bobmemecoind.service
[Unit]
Description=BobMemeCoin daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/bobmemecoin/bin/bobmemecoind -conf=/etc/bobmemecoin/bobmemecoin.conf -datadir=/var/lib/bobmemecoin
ExecStop=/opt/bobmemecoin/bin/bobmemecoin-cli -conf=/etc/bobmemecoin/bobmemecoin.conf -datadir=/var/lib/bobmemecoin stop
Restart=on-failure
RestartSec=10
TimeoutStopSec=120
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
